# rc.d/tcl-scgi.in
#!/bin/sh

# PROVIDE: tcl-scgi
# REQUIRE: LOGIN
# KEYWORD: shutdown

#
# Add the following line to /etc/rc.conf to enable tcl-scgi:
# tcl_scgi_enable="YES"
# 
# Additionally, startup flags might be defined as follows:
# tcl_scgi_flags="-addr 127.0.0.1 -port 9100"
#
# To run the server under a different user, define the following:
# tcl_scgi_user="www"
#
# The location of the scgi.tcl location might be specified as well:
# tcl_scgi_path="/usr/local/www/scgi.tcl"
#

. /etc/rc.subr

name="tcl_scgi"
rcvar=tcl_scgi_enable

load_rc_config "$name"

tcl_scgi_user=${tcl_scgi_user:-root}
tcl_scgi_path=${tcl_scgi_path:-/usr/local/www/scgi.tcl}

command="$tcl_scgi_path"
command_interpreter="@@TCLSHEXE@@"
pidfile="/var/run/tcl-scgi.pid"
start_cmd="tcl_scgi_start"
stop_cmd="tcl_scgi_stop"

tcl_scgi_start()
{
    check_startmsgs && echo -n "Starting $name"
    pid=$(su -m $tcl_scgi_user -c "$command $tcl_scgi_flags -fork")
    if [ $? -ne 0 ]; then
        return 1
    fi
    echo $pid > $pidfile

    check_startmsgs && echo "."
    return 0
}

tcl_scgi_stop()
{
    pid=$(check_pidfile $pidfile $command $command_interpreter)
    if [ -z "$pid" ]; then
        [ -n "$rc_fast" ] && return 0
        echo "$name not running? (check $pidfile)."
        return 1
    fi
    echo -n "Stopping $name"
    kill ${pid}
    rm -f $pidfile
    echo "."
}

run_rc_command "$1"
